# files/sysstat-Debian.cron
# Path to where the debian-sa1 script is located
PATH=/usr/lib/sysstat:/usr/sbin:/usr/sbin:/usr/bin:/sbin:/bin

*/1 * * * * root command -v debian-sa1 > /dev/null && debian-sa1 1 1

# Generate a daily summary of process accounting
59 23 * * * root command -v debian-sa1 > /dev/null && debian-sa1 60 2
